# products/BE-IIS-HPP-CAN-FD-SIC-REVB/test/can-gateway-testbench/can-gateway-testbench-slave.service
[Unit]
Description=CAN Gateway Testbench Slave
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/can-gateway-testbench/can_gateway_testbench.py slave --can can0 --control-bind 0.0.0.0
WorkingDirectory=/opt/can-gateway-testbench
Restart=on-failure
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target
